# scripts/openclaw-stack.service
[Unit]
Description=OpenClaw + Hermes docker compose stack
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
# CHANGE THIS to wherever you scp'd the folder (must contain docker-compose.yml + .env)
WorkingDirectory=/opt/openclaw-caddy
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
ExecReload=/usr/bin/docker compose up -d --remove-orphans
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
